# hadoop-dn/config/service/sshd/run
#!/bin/bash

echo "Starting sshd"
exec /usr/sbin/sshd -D
