# hadoop-nn-dn/config/service/serf/run
#!/bin/bash

echo "Starting serf"
exec /etc/serf/start-serf-agent.sh
